# crossword-backend/system/crossword.service
[Unit]
Description=Crossword App Service
After=network.target

[Service]
Type=simple
User=william
WorkingDirectory=/home/william/cryptic-share/crossword-backend
ExecStart=/home/william/.bun/bin/bun index.ts
Environment=PORT=3000
Environment=NODE_ENV=production
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
